# docker/sandbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:xvfb]
command=Xvfb :1 -screen 0 1280x720x24 -ac
autorestart=true
priority=100

[program:x11vnc]
command=x11vnc -display :1 -nopw -forever -shared -rfbport 5900
autorestart=true
priority=200
startsecs=2

[program:novnc]
command=websockify --web=/usr/share/novnc/ 6080 localhost:5900
autorestart=true
priority=300

[program:chromium]
command=chromium --no-sandbox --disable-gpu --disable-software-rasterizer --disable-dev-shm-usage --no-first-run --no-default-browser-check --disable-background-networking --disable-sync --disable-translate --start-maximized --window-size=1280,720 about:blank
autorestart=false
priority=400
startsecs=3
environment=DISPLAY=":1",HOME="/root"

[program:network_agent]
command=python3 /opt/network_agent.py
autorestart=true
priority=500
startsecs=5
stdout_logfile=/var/log/network_agent.log
stderr_logfile=/var/log/network_agent_err.log
environment=GATEWAY_URL="%(ENV_GATEWAY_URL)s",SANDBOX_ID="%(ENV_SANDBOX_ID)s"
